# systemd/securewatch.timer
[Unit]
Description=Run SecureWatch hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=securewatch.service

[Install]
WantedBy=timers.target
